# roles/easytier/templates/easytier-tun.service.j2
# {{ ansible_managed }}

[Unit]
Description=Configure EasyTier Interface Settings
After=easytier.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/bash {{ easytier_dir }}/easytier-tun.sh
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=multi-user.target
